# srv.service
[Unit]
Description=WebPass API server

[Service]
Type=simple
User=exedev
Group=exedev
WorkingDirectory=/home/exedev/webpass
ExecStart=/home/exedev/webpass/webpass-server -listen :8000
Restart=always
RestartSec=5
Environment=HOME=/home/exedev
Environment=USER=exedev
Environment=JWT_SECRET=CHANGE_ME_IN_PRODUCTION
Environment=CORS_ORIGINS=https://webpass.wonghome.net
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
